# systemd/srs-resolver.service
[Unit]
Description=SRS Resolver for Postfix
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/srs-resolver
Restart=on-failure
RestartSec=5s
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
